# systemd/devops-agent.service
[Unit]
Description=DevOps Agent
After=docker.service

[Service]
ExecStart=/usr/bin/devops-agent
Restart=always

[Install]
WantedBy=multi-user.target
